# service/dotfiles.timer
[Unit]
Description=Run dotfiles.service every 24 hours

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target

